# db-archive/crontab
0 2,14 * * * /backup.sh >> /backup/cron.log 2>&1
