# db-archive/crontab
0 */4 * * * /backup.sh >> /backup/cron.log 2>&1
